# deploy_tools/gunicorn-template.socket
[Unit]
Description=gunicorn-SITENAME socket

[Socket]
ListenStream=/tmp/SITENAME.socket

[Install]
WantedBy=sockets.target